# dot_config/systemd/user/battery-notifier.timer
[Unit]
Description=Battery notifier timer

[Timer]
OnBootSec=5sec
OnUnitActiveSec=1m
Unit=battery-notifier.service

[Install]
WantedBy=timers.target
