# dot_config/systemd/user/syncmail.service
[Unit]
Description=Syncmail service

[Service]
Type=oneshot
ExecStart=zsh -c -i syncmail

[Install]
WantedBy=default.target
