# dot_config/systemd/user/syncmail.timer
[Unit]
Description=Syncmail timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=5m
Unit=syncmail.service

[Install]
WantedBy=timers.target
